# docker/php/supervisord-prod.conf
# ──────────────────────────────────────────────────────────────────────────────
# Supervisord — gerencia nginx + php-fpm dentro do container de produção
# ──────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

# php-fpm inicia primeiro (nginx depende dele para repassar as requisições)
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/bin/sh -c "until nc -z 127.0.0.1 9000 2>/dev/null; do sleep 1; done && nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
